# deployment/ms1-api.service
[Unit]
Description=User & Profile Service (MS1)
After=network.target mysql.service

[Service]
Type=simple
User=ms1
WorkingDirectory=/opt/ms1-api
Environment="PATH=/opt/ms1-api/venv/bin"
Environment="DB_HOST=localhost"
Environment="DB_PORT=3306"
Environment="DB_USER=ms1_user"
Environment="DB_PASSWORD=change_me_in_production"
Environment="DB_NAME=ms1_db"
Environment="FASTAPIPORT=8000"
ExecStart=/opt/ms1-api/venv/bin/uvicorn main_db:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
